# olpc-utils-dextrose3/etc/init/olpc-log-shutdown.conf
# Upstart event
# olpc-log-shutdown - write session stop to log
#

start on starting rc RUNLEVEL=[06]

task
exec /usr/sbin/olpc-log-shutdown
